# files/etc/rc.local
#!/bin/sh 
#/usr/sbin/iptables-restore /etc/iptables/rules.v4
logger "reached rc local"
touch /tmp/rclaunch
/bin/launch-processor.sh
exit 0

